# roles/vpn_server/templates/sqvpn-server.service.j2
[Unit]
Description=sqvpn
After=network.target

[Service]
User=root
WorkingDirectory=/var/opt/sqvpn

ExecStart=sqvpn \
  exchange server.vpn-secret \
  dev quantum-vpn0 \
  listen 0.0.0.0:{{ vpn_server_port }} \
  peer client.vpn-public \
  allowed-ips {{ allowed_ips }}/{{ cidr }}

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
